# Unlockable/repo/api/supervisord.conf
[program:app]
command = python /home/dotcloud/current/api.py
